# deploy/systemd/seobuline-news-sync.timer
[Unit]
Description=Run Seobuline news sync every 30 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Unit=seobuline-news-sync.service
Persistent=true

[Install]
WantedBy=timers.target
