# deploy/systemd/seobuline-news-sync.service
[Unit]
Description=Seobuline news RSS sync job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/seobuline
ExecStart=/opt/seobuline/scripts/run-sync-news.sh
User=www-data
Group=www-data
